# crontab
# Scheduled syncs, run by Supercronic on the Fly.io `cron` Machine
# (spawned by scripts/cron-server.js — see [processes] in fly.toml).
# Deployed with the app image on every `fly deploy` — no manual scheduling
# scripts. The same Machine also serves the admin Jobs tab's on-demand
# triggers over private networking.
#
# Every job goes through scripts/run-job.js, which records the run and its
# logs in the job_runs table (visible in the admin portal's Jobs tab).
# Job definitions live in services/jobRegistry.js.
#
# All times below are UTC (Fly Machines run UTC). Afghanistan time is
# Asia/Kabul (UTC+4:30), so AFT :00 is UTC :30. Supercronic skips a run if
# the previous run of the same job is still going, and run-job.js records a
# 'skipped' row if another run of the same job is active.

# Hourly Postgres cache refresh on every AFT hour (:00), except 2:00 AM and
# 4:00 AM AFT when voice-memo-sync / classroom-sync run instead. Mirrors People,
# Ding, Applicants, Reviews, and Drive metadata, then caches voice memo audio
# (Drive → Postgres only) at the end. Classroom excluded unless
# HOURLY_CACHE_INCLUDE_CLASSROOM=true in fly.toml [env].
30 0-20,22 * * * node scripts/run-job.js hourly-cache

# Daily voice memo sync: sheet writes first (Round 2, link, date, length), then
# audio download/transcode for portal playback. 2:00 AM Afghanistan (21:30 UTC).
30 21 * * * node scripts/run-job.js voice-memo-sync

# Daily Google Classroom sync (rosters, grades, Ding change history).
# 4:00 AM Afghanistan time (23:30 UTC). Requires classroom.enabled=true.
30 23 * * * node scripts/run-job.js classroom-sync
